# examples/systemd/tnl-admin.service
# Example systemd service for tnl admin SSH tunnel
# This file is automatically created by 'sudo tnl setup' command
# Location: /etc/systemd/system/tnl-admin.service

[Unit]
Description=tnl Admin SSH Tunnel (Reverse: Gateway:20001 → Box:22)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=tunneluser
Restart=always
RestartSec=10

# Replace the following with actual values:
# - {admin_port}: Admin SSH port from gateway (e.g., 20001)
# - {ssh_port}: Gateway SSH port (usually 22)
# - {gateway_ip}: Gateway server IP or hostname

ExecStart=/usr/bin/autossh -M 0 -N \
  -o "ServerAliveInterval=30" \
  -o "ServerAliveCountMax=3" \
  -o "ExitOnForwardFailure=yes" \
  -o "StrictHostKeyChecking=no" \
  -i /home/tunneluser/.ssh/tunnel_key \
  -R {admin_port}:localhost:22 \
  -p {ssh_port} \
  tunneluser@{gateway_ip}

[Install]
WantedBy=multi-user.target
